# ansible/install/postgres_18/postgres/postgres.service.j2
# BEGIN ANSIBLE MANAGED BLOCK
[Service]
Environment=PGDATA={{ POSTGRES_DATA_DIRECTORY }}
# END ANSIBLE MANAGED BLOCK
